# crontab
0 12 * * * root python /app/src/parsing.py
